# s6/services.d/server/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec s6-setuidgid tornado /usr/local/bin/python /app/launcher.py
